# alpine-s6-consul-template-base/rootfs/etc/services.d/consul-template/run
#!/usr/bin/with-contenv sh

exec consul-template -consul ${CONSUL_DOMAIN}:${CONSUL_PORT} -config /etc/consul-template/conf.d;
